# templates/default/mon-notification.conf.erb
# Startup script for the mon_notification

description "Mon Notification daemon"
start on runlevel [2345]

console log
respawn

kill timeout 240

setgid <%= node[:mon_notification][:group] %>
setuid <%= node[:mon_notification][:user] %>
exec /usr/bin/mon-notification
